# files/etc/uci-defaults/network
#!/bin/sh

uci set network.lan.proto=dhcp
uci commit network

#  /etc/init.d/network restart # Needed?

exit 0
